# deploy/supervisor.conf
; supervisor.conf — DialFlow Pro
; Manages all DialFlow processes. Copy to /etc/supervisor/conf.d/dialflow.conf
; Then: supervisorctl reread && supervisorctl update

[group:dialflow]
programs=dialflow_web,dialflow_worker,dialflow_beat

; ── Django / Daphne (ASGI) ───────────────────────────────────────────────────
[program:dialflow_web]
command=/srv/dialflow/venv/bin/daphne -b 0.0.0.0 -p 8000 dialflow.asgi:application
directory=/srv/dialflow
user=dialflow
environment=DJANGO_SETTINGS_MODULE="dialflow.settings.prod"
autostart=true
autorestart=true
stdout_logfile=/var/log/dialflow/web.log
stderr_logfile=/var/log/dialflow/web.err
stopwaitsecs=10

; ── Celery Worker ─────────────────────────────────────────────────────────────
[program:dialflow_worker]
command=/srv/dialflow/venv/bin/celery -A dialflow worker
    -l info
    --concurrency=4
    -Q celery,default
directory=/srv/dialflow
user=dialflow
environment=DJANGO_SETTINGS_MODULE="dialflow.settings.prod"
autostart=true
autorestart=true
stopwaitsecs=30
stdout_logfile=/var/log/dialflow/worker.log
stderr_logfile=/var/log/dialflow/worker.err

; ── Celery Beat ───────────────────────────────────────────────────────────────
[program:dialflow_beat]
command=/srv/dialflow/venv/bin/celery -A dialflow beat
    -l info
    --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/srv/dialflow
user=dialflow
environment=DJANGO_SETTINGS_MODULE="dialflow.settings.prod"
autostart=true
autorestart=true
stdout_logfile=/var/log/dialflow/beat.log
stderr_logfile=/var/log/dialflow/beat.err

; Note: The ARI worker starts automatically inside Daphne via TelephonyConfig.ready()
; No separate process needed.
